# dockstuff/supervisord.conf
[supervisord]
nodaemon=false

[program:stabman]
command=/client $HADOOP_HOME
